# meross_local_broker/rootfs/etc/services.d/mosquitto/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree based on service exit code
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

# Unregister mosquitto service
bashio::log.info "Unregistering the Mosquitto daemon"
mkdir -p /var/run/mdns/meross-mqtt
path=$(cat /var/run/mdns/meross-mqtt/service_path)
/opt/custom_broker/mdns_register.py unregister $path
rm -R path

/run/s6/basedir/bin/halt