# meross_local_broker/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start meross API service
# ==============================================================================
# Advertising api_port service:
bashio::log.info "Avertising the Meross API service"
api_port=$(bashio::addon.port 2002)
mkdir -p /var/run/mdns/meross-http-api
/opt/custom_broker/mdns_register.py register meross-local-api meross-local-api _meross-mqtt._tcp $api_port --update-if-present > /var/run/mdns/meross-http-api/service_path

# Starting the API service:
bashio::log.debug "Starting the Meross daemon SERVICE"
exec 2>&1
exec /usr/local/bin/api.sh
